# home/config/xcape/.config/autostart/xcape.desktop
[Desktop Entry]
Type=Application
Name=xcape
Exec=/home/andrea/bin/xcape.sh
Terminal=false
